# teaching_panel.service
[Unit]
Description=Teaching Panel Gunicorn Service (flatten-ready)
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
# Дополнительные каталоги и PID (runtime isolation)
RuntimeDirectory=teaching_panel
PIDFile=/run/teaching_panel/teaching_panel.pid
# Текущая вложенная структура: manage.py находится в /var/www/teaching_panel/teaching_panel/teaching_panel
# Если позже распакуете структуру (уберёте один уровень), обновите WorkingDirectory.
WorkingDirectory=/var/www/teaching_panel/teaching_panel/teaching_panel
Environment="PATH=/var/www/teaching_panel/venv/bin"
EnvironmentFile=-/var/www/teaching_panel/teaching_panel/teaching_panel/.env
# Frontend URL for payment redirects
Environment="FRONTEND_URL=http://72.56.81.163"
# Production settings
Environment="ALLOWED_HOSTS=72.56.81.163,localhost,127.0.0.1"
# Environment="DJANGO_SETTINGS_MODULE=teaching_panel.settings"
# Environment="SECRET_KEY=__CHANGE_ME__"
# Environment="DEBUG=False"
Environment="GUNICORN_WORKERS=5"
Environment="GUNICORN_BIND=0.0.0.0:8000"
Environment="GUNICORN_TIMEOUT=120"

# Убрали лишний --chdir: WorkingDirectory уже задаёт нужный путь
ExecStart=/var/www/teaching_panel/venv/bin/gunicorn teaching_panel.wsgi:application \
  --bind ${GUNICORN_BIND} \
  --workers ${GUNICORN_WORKERS} \
  --timeout ${GUNICORN_TIMEOUT} \
  --graceful-timeout ${GUNICORN_TIMEOUT} \
  --max-requests 1000 \
  --max-requests-jitter 100 \
  --keep-alive 5 \
  --log-level info \
  --access-logfile - \
  --error-logfile -

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=10
PrivateTmp=true
Restart=always
StartLimitIntervalSec=300
StartLimitBurst=5
RestartSec=3

[Install]
WantedBy=multi-user.target

# DEPLOY СТЕПЫ:
# 1) scp teaching_panel.service root@SERVER:/etc/systemd/system/teaching_panel.service
# 2) systemctl daemon-reload
# 3) systemctl enable teaching_panel --now
# 4) journalctl -u teaching_panel -n 50 --no-pager
# 5) curl -X POST http://127.0.0.1:8000/api/jwt/token/ -H 'Content-Type: application/json' -d '{"email":"user@example.com","password":"Pass123"}'

# FLATTEN (опционально позже):
# mv /var/www/teaching_panel/teaching_panel/teaching_panel/* /var/www/teaching_panel/teaching_panel/
# mv /var/www/teaching_panel/teaching_panel/teaching_panel/.* /var/www/teaching_panel/teaching_panel/ 2>/dev/null || true
# rmdir /var/www/teaching_panel/teaching_panel/teaching_panel
# Обновить WorkingDirectory и ExecStart при необходимости.
