# deploy/teaching_panel_stage.service
[Unit]
Description=Teaching Panel Staging - Gunicorn Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
RuntimeDirectory=teaching_panel_stage
PIDFile=/run/teaching_panel_stage/teaching_panel_stage.pid
WorkingDirectory=/var/www/teaching-panel-stage/teaching_panel
Environment="PATH=/var/www/teaching-panel-stage/venv/bin"
EnvironmentFile=-/var/www/teaching-panel-stage/.env
Environment="DJANGO_SETTINGS_MODULE=teaching_panel.settings_staging"
Environment="FRONTEND_URL=https://stage.lectiospace.ru"
Environment="ALLOWED_HOSTS=stage.lectiospace.ru,localhost,127.0.0.1"
Environment="GUNICORN_WORKERS=2"
Environment="GUNICORN_BIND=127.0.0.1:8001"
Environment="GUNICORN_TIMEOUT=120"

ExecStart=/var/www/teaching-panel-stage/venv/bin/gunicorn teaching_panel.wsgi:application \
  --bind ${GUNICORN_BIND} \
  --workers ${GUNICORN_WORKERS} \
  --timeout ${GUNICORN_TIMEOUT} \
  --graceful-timeout ${GUNICORN_TIMEOUT} \
  --max-requests 1000 \
  --max-requests-jitter 100 \
  --keep-alive 5 \
  --log-level info \
  --access-logfile - \
  --error-logfile -

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=10
PrivateTmp=true
Restart=always
StartLimitIntervalSec=300
StartLimitBurst=5
RestartSec=3

[Install]
WantedBy=multi-user.target

# DEPLOY:
# 1) scp deploy/teaching_panel_stage.service root@SERVER:/etc/systemd/system/teaching_panel_stage.service
# 2) systemctl daemon-reload
# 3) systemctl enable teaching_panel_stage --now
# 4) journalctl -u teaching_panel_stage -n 50 --no-pager
