# deploy/temp_teaching_panel_stage_b.service
# =============================================================================
# STAGE B: Growth Configuration (750 teachers / 7,500 students)
# Server: 4 vCPU / 8 GB RAM
# Expected: 120-180 RPS peak
# Applied: 2026-02-02
# =============================================================================

[Unit]
Description=Teaching Panel Django Application (Stage B - Growth)
After=network.target postgresql.service
Wants=nginx.service

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/var/www/teaching_panel/teaching_panel
Environment="DJANGO_SETTINGS_MODULE=teaching_panel.settings"
Environment="PATH=/var/www/teaching_panel/venv/bin"
EnvironmentFile=-/var/www/teaching_panel/.env

# STAGE B CONFIGURATION:
# - 5 gevent workers (optimal for 4 vCPU)
# - gevent for async I/O (handles DB waits efficiently)
# - 1000 connections per worker
# - 60s timeout (balanced)
ExecStart=/var/www/teaching_panel/venv/bin/gunicorn \
    --bind 127.0.0.1:8000 \
    --workers 5 \
    --worker-class gevent \
    --worker-connections 1000 \
    --timeout 60 \
    --graceful-timeout 30 \
    --keep-alive 5 \
    --max-requests 2000 \
    --max-requests-jitter 200 \
    --access-logfile /var/log/teaching_panel/access.log \
    --error-logfile /var/log/teaching_panel/error.log \
    --capture-output \
    --preload \
    teaching_panel.wsgi:application

Restart=always
RestartSec=3
StartLimitIntervalSec=300
StartLimitBurst=5

# Resource limits for 8GB server
# Workers: 5 × ~120MB = 600MB + buffer
MemoryMax=1G
MemoryHigh=800M

CPUWeight=100

[Install]
WantedBy=multi-user.target
