# deploy/scaling/stage-a/celery-beat.service
# =============================================================================
# STAGE A: Celery Beat Scheduler
# =============================================================================

[Unit]
Description=Teaching Panel Celery Beat (Stage A)
After=network.target redis.service
Wants=celery.service

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/teaching_panel/teaching_panel
Environment="PATH=/var/www/teaching_panel/venv/bin"
EnvironmentFile=-/var/www/teaching_panel/.env

ExecStart=/var/www/teaching_panel/venv/bin/celery -A teaching_panel beat \
    --loglevel=info \
    --scheduler django_celery_beat.schedulers:DatabaseScheduler

# Beat is lightweight - minimal memory
MemoryMax=150M

Restart=always
RestartSec=10

[Install]
WantedBy=multi-user.target
