# teaching_panel/deployment/teaching_panel.service
[Unit]
Description=Teaching Panel Gunicorn Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/teaching_panel
Environment="PATH=/var/www/teaching_panel/venv/bin"
EnvironmentFile=-/var/www/teaching_panel/.env
# (Optional) '-' before EnvironmentFile: missing .env non-fatal. settings.py loads /var/www/teaching_panel/.env if present.
# IMPORTANT: Keep ExecStart on a single line (avoid CRLF + line continuation issues under systemd).
ExecStart=/var/www/teaching_panel/venv/bin/gunicorn --bind 0.0.0.0:8000 --workers 3 --worker-class sync --timeout 120 --graceful-timeout 120 --max-requests 1000 --max-requests-jitter 100 --log-level info --access-logfile - --error-logfile - teaching_panel.wsgi:application

ExecReload=/bin/kill -s HUP $MAINPID
KillMode=mixed
TimeoutStopSec=10
PrivateTmp=true
Restart=always

[Install]
WantedBy=multi-user.target
